# ambari/configs/supervisord.d/cluster-deploy.conf
[program:cluster-deploy]
command=/cluster-deploy.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=false
priority=3
